# cust/root/etc/rc.local
#!/bin/sh -eu

/usr/local/lib/rescuedisk/setup.sh
